# deploy/supervisor/hawkeye.conf
[program:hawkeye]
command=/usr/local/bin/gunicorn -w10 -b127.0.0.1:8888 api:app
directory=/Hawkeye/server
startsecs=5
stopwaitsecs=0
autostart=true
autorestart=true