# deploy/supervisor/huey.conf
[program:huey]
command=/usr/local/bin/huey_consumer.py task.huey -w 5
directory=/Hawkeye/server
startsecs=5
stopwaitsecs=0
autostart=true
autorestart=true