# deploy/supervisor/huey.conf
[program:huey]
command=/usr/local/bin/huey_consumer.py task.huey -k process -w 4
directory=/Hawkeye/server
startsecs=5
stopwaitsecs=0
autostart=true
autorestart=true
