# deploy/supervisor/openresty.conf
[program:openresty]
command=/usr/local/openresty/bin/openresty -g "daemon off;"
directory=/tmp
startsecs=5
stopwaitsecs=0
autostart=true
autorestart=true