# deploy/supervisor/redis.conf
[program:redis]
command=/usr/bin/redis-server
directory=/tmp
startsecs=10
stopwaitsecs=0
autostart=true
autorestart=true